# infra/systemd/sunyata-ai.service
[Unit]
Description=Sunyata AI Microservice (FastAPI)
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/var/www/sunyata/services/ai
Environment="PATH=/var/www/sunyata/services/ai/.venv/bin:/usr/local/bin:/usr/bin"
EnvironmentFile=/var/www/sunyata/services/ai/.env
ExecStart=/var/www/sunyata/services/ai/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/www/sunyata/services/ai/logs
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
